# configs/supervisor.blog.conf
[program:blog]
directory=/var/www/blog
command=/usr/local/bin/gunicorn wsgi -c /var/www/blog/configs/gunicorn.blog.py
autostart=true
autorestart=true
