# configs/supervisor.blog.conf
[program:blog]
directory=/var/www/blog/configs
command=/usr/local/bin/gunicorn /var/www/blog/wsgi -c gunicorn.blog.py
autostart=true
autorestart=true
